# cronjob_gentesexy.cron
0 23 * * * /Users/josegalarza/Code/Python/gentesexy/cronjob_gentesexy.sh
